# host-server/etc/cron.d/kindle-weather
0,30 * * * * root sh -c "/var/lib/kindle-weather-host/kindle-weather.sh"
